# deploy/supervisor/webadmin-celeryd.conf
#Copyright (c) 2015-present, Philippine-California Advanced Research Institutes-
#The Village Base Station Project (PCARI-VBTS). All rights reserved.
#
#This source code is licensed under the BSD-style license found in the
#LICENSE file in the root directory of this source tree.

[program:webadmin-celeryd]

directory=/var/www/vbts_webadmin
command=envdir /var/www/vbts_webadmin/envdir celery -A vbts_webadmin worker -l info --statedb=worker.state
user=www-data
group=www-data
stdout_logfile=/var/log/webadmin-celeryd.log
stderr_logfile=/var/log/webadmin-celeryd.log
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 60
killasgroup=true
